# v6_IBKR_WebAPI/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:botpy]
command=bash -c "python gateway/wait_for_gateway.py --port %(ENV_IBKR_PORT)s --timeout 300 && python main.py"
autostart=true
autorestart=true
startsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
